# infra/systemd/ib-altdata-qa.service
[Unit]
Description=IB Bot daily versioned QA for altdata_snapshots
After=docker.service network-online.target
Wants=network-online.target
OnFailure=ib-altdata-qa-alert.service

[Service]
Type=oneshot
User=servidor
Group=servidor
WorkingDirectory=/home/servidor/Desktop/cursor-projects/ib_bot
ExecStart=/home/servidor/Desktop/cursor-projects/ib_bot/infra/scripts/run_altdata_qa_daily.sh
TimeoutStartSec=900
Slice=workloads.slice
MemoryMax=1G
CPUQuota=200%
IOSchedulingClass=idle
Nice=10
UMask=0002
