# infra/systemd/ib-backtests.service.d/10-price-source.conf
# Phase f2_backtest_fix_alerta (Finding:5): the weekly backtest must run with the
# IB Gateway OFF (phase f1 keeps ibgateway/ibeam dormant). Force the price source
# to yfinance so neither the 56 strategy backtests nor the dashboard step touch
# live IB.
#
# NOTE: this is passed as an explicit --price-source ARGUMENT, not via
# Environment=. The unit's EnvironmentFile=.env pins PRICE_SOURCE=ib and, by
# systemd's env-merge order, overrides a plain Environment= override — but
# run_all_backtests.py resolves `args.price_source or os.getenv("PRICE_SOURCE")`,
# so the CLI arg wins unconditionally. yfinance's api_caution budget is
# warn=10k/block=60k (vs ib warn=1.5k) so the run is never refused; the dashboard
# step itself runs with --cache-only (see run_all_backtests.py).
[Service]
ExecStart=
ExecStart=/home/servidor/Desktop/cursor-projects/ib_bot/.venv/bin/python run_all_backtests.py --price-source yfinance
Environment=PRICE_SOURCE=yfinance
