# infra/systemd/ib-backtests.service.d/10-price-source.conf
# Phase f2_backtest_fix_alerta (Finding:5): the weekly backtest must run with the
# IB Gateway OFF (phase f1 keeps ibgateway/ibeam dormant). Force the price source
# to yfinance so neither the 56 strategy backtests nor the dashboard step touch
# live IB. yfinance's api_caution budget is warn=10k/block=60k, so the run is not
# refused; the dashboard step itself runs with --cache-only (see run_all_backtests.py).
[Service]
Environment=PRICE_SOURCE=yfinance
