# infra/systemd/ib-backtests.service.d/15-preflight-containment.conf
# Keep the hour-long weekly run inside the EPYC workload envelope and fail
# before doing any backtest work if the canonical dashboard generator is absent.
# The phase-f2 worktree override was intentionally temporary; production always
# runs the merged files from the main checkout.
[Service]
ExecStartPre=/usr/bin/test -f /home/servidor/Desktop/cursor-projects/ib_bot/run_all_backtests.py
ExecStartPre=/usr/bin/test -f /home/servidor/Desktop/cursor-projects/ib_bot/generate_plot_data.py
Slice=workloads.slice
MemoryMax=24G
MemoryHigh=24G
# Bound the process-local yfinance LRU. The disk cache is canonical; retaining
# all ~2,700 price frames made a warm weekly run needlessly reclaim-heavy.
Environment=YF_MEMORY_CACHE_MAX_TICKERS=128
CPUQuota=2200%
CPUWeight=30
IOWeight=50
IOSchedulingClass=idle
Nice=10
# earlyoom already prefers Python compute processes. Keep the neutral kernel
# score so victim choice is proportional to actual memory use; +800 caused this
# ~250 MiB job to be terminated before a co-running ~49 GiB workload.
OOMScoreAdjust=0
# systemd-oomd does not use OOMScoreAdjust. The weekly backtest is a small,
# fail-closed maintenance job whose cache-heavy access pattern can show more
# reclaim pressure than a much larger compute job. Deprioritize it as an oomd
# victim while retaining MemoryHigh/MemoryMax as the hard safety boundary.
ManagedOOMPreference=avoid
