# infra/systemd/ib-backtests.service.d/15-preflight-containment.conf
# Keep the hour-long weekly run inside the EPYC workload envelope and fail
# before doing any backtest work if the canonical dashboard generator is absent.
# The phase-f2 worktree override was intentionally temporary; production always
# runs the merged files from the main checkout.
[Service]
ExecStartPre=/usr/bin/test -f /home/servidor/Desktop/cursor-projects/ib_bot/run_all_backtests.py
ExecStartPre=/usr/bin/test -f /home/servidor/Desktop/cursor-projects/ib_bot/generate_plot_data.py
Slice=workloads.slice
MemoryMax=24G
MemoryHigh=24G
CPUQuota=2200%
CPUWeight=30
IOWeight=50
IOSchedulingClass=idle
Nice=10
OOMScoreAdjust=800
