# infra/systemd/ib-backtests.service.d/20-onfailure.conf
# Phase f2_backtest_fix_alerta (Finding:5): on ANY failed run, automatically wake
# the ib_bot Domain Manager in Conductor so it resolves the failure autonomously.
[Unit]
OnFailure=ib-backtests-alert.service
